# x735fan.service
[Unit]
Description=X735 Fan Speed Control

[Service]
Type=simple
ExecStart=scriptpath/x735fan.sh

[Install]
WantedBy=multi-user.target
